# Projects/IHC/IHC_Services/IHC_Controller/IHC_Controller_Service/esphome_server.service
[Unit]
Description=ESPHome_Server

[Service]
ExecStart=/home/pi/ESPHome/ESPHome_Server
WorkingDirectory=/home/pi/ESPHome
Restart=always

[Install]
WantedBy=multi-user.target
